# config/com.user.iphone-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for scheduled iPhone backups.

  Setup (handled automatically by install.sh):
    1. Build the binary:  cargo build -\-release
    2. Install it:        cp target/release/iphone-backup /usr/local/bin/
    3. Install this file: cp config/com.user.iphone-backup.plist ~/Library/LaunchAgents/
    4. Load:              launchctl load ~/Library/LaunchAgents/com.user.iphone-backup.plist

  Manual trigger:
    launchctl start com.user.iphone-backup

  Disable:
    launchctl unload ~/Library/LaunchAgents/com.user.iphone-backup.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.iphone-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/iphone-backup</string>
        <string>backup</string>
    </array>

    <!-- Run daily at 2:00 am -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/iphone-backup-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/iphone-backup-launchd.log</string>
</dict>
</plist>
